# site-modules/consul_template/templates/consul-template.service.erb
# This file is managed by Puppet.
[Unit]
Description=consul-template
Requires=network-online.target
After=network-online.target consul.service vault.service

[Service]
EnvironmentFile=-/etc/sysconfig/consul-template
Restart=on-failure
ExecStart=<%= scope['consul_template::install_dir'] %>/<%= scope['consul_template::bin'] %> -config=<%= scope['consul_template::config_dir'] %>/config.hcl

[Install]
WantedBy=multi-user.target
